# systemd/github-pr-review.service
[Unit]
Description=GitHub PR Review - Review random PRs

[Service]
Type=oneshot
ExecStart=/bin/bash /home/yhm/.openclaw/scripts/github-pr-review.sh
User=yhm
Environment=PATH=/home/yhm/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
